# asterisk/taxi-bridge.service
[Unit]
Description=Taxi AI Asterisk Bridge
Documentation=https://github.com/your-org/taxi-ai
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
User=asterisk
Group=asterisk
WorkingDirectory=/opt/taxi-bridge
ExecStart=/usr/bin/python3 -u /opt/taxi-bridge/taxi_bridge_v7.py
Restart=always
RestartSec=5

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=WS_URL=wss://oerketnvlmptpfvttysy.supabase.co/functions/v1/taxi-realtime-paired

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taxi-bridge

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target
